# packaging/systemd/oisp-sensor.service
[Unit]
Description=OISP Sensor - Universal AI Observability
Documentation=https://sensor.oisp.dev
After=network.target

[Service]
Type=simple
User=root
Group=root

# Main command - records AI activity with web UI on default port
ExecStart=/usr/local/bin/oisp-sensor record \
    --config /etc/oisp-sensor/config.toml \
    --output /var/log/oisp-sensor/events.jsonl \
    --port 7777

# Reload configuration on SIGHUP
ExecReload=/bin/kill -HUP $MAINPID

# Restart on failure with exponential backoff
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=SIGTERM

# Stop timeout
TimeoutStopSec=30

# Create required directories before starting
ExecStartPre=/bin/mkdir -p /var/log/oisp-sensor
ExecStartPre=/bin/mkdir -p /var/lib/oisp-sensor

# Security hardening (compatible with eBPF requirements)
NoNewPrivileges=no
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/oisp-sensor /var/lib/oisp-sensor
PrivateTmp=true

# Required capabilities for eBPF SSL capture
# CAP_SYS_ADMIN: Required for BPF operations on older kernels
# CAP_BPF: Required for loading BPF programs (kernel >= 5.8)
# CAP_PERFMON: Required for perf events (kernel >= 5.8)
# CAP_NET_ADMIN: Required for network-related BPF operations
AmbientCapabilities=CAP_SYS_ADMIN CAP_BPF CAP_PERFMON CAP_NET_ADMIN
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_BPF CAP_PERFMON CAP_NET_ADMIN

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oisp-sensor

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Resource limits
LimitMEMLOCK=infinity
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

